# scripts/com.username.memory-sync.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.YOUR-USERNAME.memory-sync</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/ABSOLUTE/PATH/TO/memory_sync.py</string>
    </array>

    <!-- Run every 5 minutes (300 seconds). Change StartInterval to adjust. -->
    <key>StartInterval</key>
    <integer>300</integer>

    <!-- Also run immediately after login -->
    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/YOUR-USERNAME/.gemini/antigravity-ide/sync_memory_to_ki.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/YOUR-USERNAME/.gemini/antigravity-ide/sync_memory_to_ki_err.log</string>
</dict>
</plist>
